# nautilus-mission.service
# nautilus-mission.service - autonomous hard-code qualify on boot (NO tmux).
#
# Runs nautilus_autostart.sh directly, which brings up the stack, waits for it
# to be ready, then runs the mission. Output -> journalctl -u nautilus-mission.
#
# TESTBENCH <-> LIVE toggle:
#     sudo systemctl enable  nautilus-mission     # LIVE: runs on every boot
#     sudo systemctl disable nautilus-mission     # TESTBENCH: boots do nothing
# Watch it:
#     journalctl -u nautilus-mission -f
# Stop a running mission:
#     sudo systemctl stop nautilus-mission
#
# INSTALL:
#     sudo cp ~/nautilus_ws/scripts/nautilus-mission.service /etc/systemd/system/
#     sudo systemctl daemon-reload
#
# *** The HARDWARE KILL SWITCH is the interlock. With this enabled, every boot
# arms and drives once the stack is up. Keep it engaged for any bench boot. ***

[Unit]
Description=Nautilus autonomous hard-code qualify mission
After=multi-user.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nano
ExecStart=/home/nano/nautilus_ws/scripts/nautilus_autostart.sh
# stack bring-up + countdown can take a while; don't let systemd time it out
TimeoutStartSec=0
# finite mission -- run once per boot, don't respawn in a loop
Restart=no

[Install]
WantedBy=multi-user.target
